# camgrab.service
[Unit]
Description=Camgrab Script

[Service]
ExecStart=/usr/local/bin/camgrab
User=camera

[Install]
WantedBy=multi-user.target
